# systemd/arc-dispatch.service
[Unit]
Description=arc-agent dispatch runner
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/dev/arc0btc/arc-agent
ExecStart=/usr/bin/env bun src/cli.ts run
Environment="HOME=%h"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:%h/.bun/bin"
TimeoutStopSec=3600
StandardOutput=journal
StandardError=journal
