# systemd/arc-dispatch.timer
[Unit]
Description=arc-agent dispatch timer — fires every 1 minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
